# crontab
0 0 1 * * cd /usr/app && /usr/local/bin/node /usr/app/index.js >> /var/log/cron.log 2>&1